# provisioning/hayward-provisioning.service
[Unit]
Description=Hayward HeatPro WiFi AP Provisioning
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStartPre=/usr/bin/env mkdir -p /var/run/hayward
ExecStart=/usr/local/bin/hayward-provisioning-server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
